# deploy/releaselog-ingest.service
[Unit]
Description=ReleaseLog candidate ingestion (oneshot)
After=releaselog.service
Requires=releaselog.service

[Service]
Type=oneshot
ExecStart=/root/releaselog/deploy/call-cron.sh /api/cron/ingest
User=root
