# deploy/releaselog-weekly-digest.timer
[Unit]
Description=Run ReleaseLog weekly digest Mondays 09:00 UTC

[Timer]
OnCalendar=Mon *-*-* 09:00:00 UTC
Persistent=true
Unit=releaselog-weekly-digest.service

[Install]
WantedBy=timers.target
