# deploy/releaselog-ingest.timer
[Unit]
Description=Run ReleaseLog candidate ingestion every 30 minutes

[Timer]
OnCalendar=*:0/30
Persistent=true
Unit=releaselog-ingest.service

[Install]
WantedBy=timers.target
