# deploy/releaselog-send.service
[Unit]
Description=ReleaseLog notification send worker (oneshot)
After=releaselog.service
Requires=releaselog.service

[Service]
Type=oneshot
ExecStart=/root/releaselog/deploy/call-cron.sh /api/cron/send-notifications
User=root
